# packaging/linux-agent-web.service
# Linux 运维 Agent — Web 控制台 systemd 单元（生产加固模板）
#
# 用途：把 Web 控制台作为长期服务运行在目标主机上，以专用非 root 用户执行，
# 并施加 systemd 沙箱与 cgroup 资源约束。远程访问一律通过 SSH 端口转发到
# 127.0.0.1:8765，不直接暴露到 LAN/公网。
#
# 安装步骤：下载固定版本的 linux-agent-install.sh 后执行：
#   sudo bash linux-agent-install.sh install --version vX.Y.Z --require-signature
# 首次安装会临时启动完成健康检查，随后停止服务且不改启用状态；长期运行需显式执行：
#   sudo systemctl enable --now linux-agent-observer-helper.socket linux-agent-web.service
# 升级、回滚和健康检查分别使用 upgrade、rollback、health 子命令。
#
# 访问：ssh -L 8765:127.0.0.1:8765 <user>@<target-host>，然后浏览器打开 127.0.0.1:8765。

[Unit]
Description=Linux 运维 Agent Web 控制台
Documentation=file:/opt/linux-agent/README.md
After=network-online.target
Wants=network-online.target
Requires=linux-agent-observer-helper.socket
After=linux-agent-observer-helper.socket

[Service]
Type=notify
# —— 身份边界：专用非 root 账户，禁止提权 ——
User=linux-agent
Group=linux-agent
WorkingDirectory=/opt/linux-agent/current
ExecStart=/usr/bin/env bash /opt/linux-agent/current/bin/agent-web
Restart=on-failure
RestartSec=3s

# —— 进程/执行边界 ——
NoNewPrivileges=yes
# 该服务本身不需要提权；如需 auditd observer 提权，改由独立 privileged helper 承担。

# —— 文件系统边界 ——
ProtectSystem=strict
ProtectHome=yes
# 代码、Provider registry、schema 和提示词保持只读；仅开放运行时数据目录。
ReadWritePaths=/opt/linux-agent/data
PrivateTmp=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectProc=invisible
ProcSubset=pid
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @swap @reboot @raw-io

# —— 网络出口边界 ——
# Web 只监听 127.0.0.1（由 bin/agent-web 在 remote 模式强制）。
# 生产环境按 packaging/dropins/10-provider-egress.conf.example 启用默认拒绝，
# 再显式加入 Provider CIDR；域名解析与动态地址变更由运维侧维护。

# —— cgroup 资源硬上限 ——
CPUQuota=200%
MemoryMax=1G
TasksMax=256
# 单个服务实例的墙钟上限交由应用层 execution.timeout_sec 控制；
# 若需服务级重启周期，可加 RuntimeMaxSec=。

[Install]
WantedBy=multi-user.target
